# imanager/imanager.service.template
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=$service_start

[Install]
WantedBy=multi-user.target
