# lion_canbus/canbus.service.template
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=$can_init

[Install]
WantedBy=multi-user.target
EOF
